# systemd/telnet-proxy@.service
[Unit]
Description=Telnet Proxy
After=local-fs.target

[Service]
ExecStart=-/usr/local/bin/telnet-proxy
User=<user>
StandardInput=socket
