# systemd/telnet-proxy.socket
[Unit]
Description=Telnet Proxy

[Socket]
ListenStream=<port>
Accept=true
MaxConnections=3

[Install]
WantedBy=sockets.target
